# supervisord.conf
# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:gunicorn]
command=gunicorn --workers 3 --bind 127.0.0.1:5000 app:app
directory=/app
autostart=true
autorestart=true